# chore-tracker-addon/rootfs/etc/services.d/chore-tracker/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info "Starting Chore Tracker API..."
exec python3 /usr/bin/chore-tracker
